# src/resources/wakka-shutdown.service
[Unit]
Description=Wakka — Install pending updates before shutdown
# Run before shutdown/reboot/halt targets
DefaultDependencies=no
Before=poweroff.target reboot.target halt.target shutdown.target
Conflicts=rescue.service rescue.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/wakka-shutdown-helper
# Give updates up to 1 hour to complete
TimeoutStartSec=3600
RemainAfterExit=no

[Install]
WantedBy=poweroff.target reboot.target halt.target
