# main/install/wakka-shutdown.service
[Unit]
Description=Wakka — Install pending updates before shutdown
# Run before shutdown/reboot/halt targets
DefaultDependencies=no
Before=poweroff.target reboot.target halt.target shutdown.target
Conflicts=rescue.service rescue.target
After=network-online.target graphical.target

[Service]
Type=oneshot
# Run as the first non-root user (adjust if needed, or use User= with the actual user)
# The helper is designed to be run as the user, yay handles sudo internally
ExecStart=/usr/bin/wakka-shutdown-helper
# Give updates up to 1 hour to complete
TimeoutStartSec=3600
# Keep running until explicitly stopped
RemainAfterExit=no
# Environment
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus

[Install]
WantedBy=poweroff.target reboot.target halt.target
